# .config/systemd/user.control/mail.service
[Unit]
Description=Check Email

[Service]
ExecStart=/usr/bin/mbsync --all --verbose
