# .config/systemd/user.control/vdirsyncer.service
[Unit]
Description=Sync Contacts and Calendars

[Service]
ExecStart=/usr/bin/vdirsyncer metasync
