# .config/systemd/user/vdirsyncer.timer
[Unit]
Description="Scheduled calendar and contact sync"

[Timer]
OnCalendar="1h"

[Install]
WantedBy=default.target
