# .config/systemd/user/vdirsyncer.service
[Unit]
Description = Sync contacts and calendar

[Service]
Type = oneshot
ExecStart = vdirsyncer sync
