# .config/systemd/user/vdirsyncer_metasync.service
[Unit]
Description = Sync contacts and calendar metadata

[Service]
Type = oneshot
ExecStart = vdirsyncer metasync
