# .config/systemd/user/vdirsyncer.timer
[Unit]
Description = Sync contacts and calendar

[Timer]
OnCalendar = Hourly

[Install]
WantedBy = timers.target
